# 80-bt-headset.rules
ACTION=="add", SUBSYSTEM=="bluetooth", RUN+="/home/baizon/GitHub/a2dpfix/a2dp-fix-wrapper"
